# deploy/initiales-backend.service
# /etc/systemd/system/initiales-backend.service
# À installer sur le VPS.

[Unit]
Description=Initiales Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/initiales/backend
Environment="PATH=/var/www/initiales/backend/venv/bin"
EnvironmentFile=/var/www/initiales/backend/.env
ExecStart=/var/www/initiales/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/initiales

[Install]
WantedBy=multi-user.target

# Installation :
# sudo cp deploy/initiales-backend.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable --now initiales-backend
# sudo systemctl status initiales-backend
# sudo journalctl -u initiales-backend -f
